# rpm_install/etc/cron.d/ckan-harvest
*/15 * * * * root ckan --plugin=ckanext-harvest harvester run -c /etc/ckan/production.ini >> /var/log/harvester_run.log 2>&1
*/15 * * * * root ckan --plugin=ckanext-geoserver geoserver publish-ogc-redis-queue -c /etc/ckan/production.ini >> /var/log/publish_ogc_queue.log 2>&1
*/5 * * * * root ckan --plugin=ckanext-ngds  ngds purge-deleted-harvest-sources -c /etc/ckan/production.ini >> /var/log/purge_deleted_harvest_sources.log 3>&1
